# archlive/airootfs/etc/skel/.config/autostart/archinstall_gnome-terminal.desktop
[Desktop Entry]
Type=Application
Exec=gnome-terminal -e sudo archinstall
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true