# archlive/airootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=/usr/bin/systemctl enable --now gdm.service
ExecStart=-/sbin/agetty -o '-p -f -- \\u' --noclear --autologin liveuser - $TERM
